# keystone/keystone-2/cron.d/keystone-fernet-rotate
#Ansible: Fernet auto rotate job
@daily keystone /opt/keystone-fernet-rotate.sh
