# keystone/keystone-2/cron.d/keystone-credential-rotate
#Ansible: Credential auto rotate job
@weekly keystone /opt/keystone-credential-rotate.sh
